# packaging/systemd/user/avreamd.service
[Unit]
Description=AVream daemon (user)
After=default.target

[Service]
Type=simple
EnvironmentFile=-%h/.config/avream/avreamd.env
ExecStart=avreamd
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
